# sentinel/files/sentinel.service
[Unit]
Description=Sentinel — Creator Analytics health check
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/sentinel
ExecStart=/usr/bin/python3 /opt/sentinel/health_check.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
